# config/systemd/user/inhibit-suspend-on-ssh.service
[Unit]
Description=Inhibit suspend when SSH connection is active
After=network.target sshd.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.config/sway/scripts/inhibit-suspend-on-ssh.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
