# config/systemd/user/set-bing-wallpaper.service
[Unit]
Description=Set Bing Wallpaper

[Service]
ExecStart=%h/.config/sway/set-bing-wallpaper.sh
